# daemons/daemon-updater.ini
[program:daemon-updater]
command=python %(ENV_SUPERVISOR_CONFIG_REPO_LOCN)s/webhook/example.py
redirect_stderr=true
stdout_logfile=/tmp/daemon-updater.log
